# hw7.conf
description "uWSGI server instance configured to serve hw7"

start on runlevel [2345]
stop on runlevel [!2345]

setuid root
setgid www-data

chdir /root/hw7/
exec uwsgi --ini hw7.ini
